# bluetooth-agent.service
[Unit]
Description=IEEE Presence Tracker Bluetooth Agent
Documentation=https://github.com/ieee/presence-tracker
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStart=/home/ieee/.local/bin/uv run --project "/home/ieee/Desktop/IEEE Presence Tracker" "/home/ieee/Desktop/IEEE Presence Tracker/bluetooth_agent.py"
WorkingDirectory=/home/ieee/Desktop/IEEE Presence Tracker
User=root
Group=bluetooth
Restart=always
RestartSec=5

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
