# bt-agent.service
[Unit]
Description=Bluetooth Pairing Agent
Documentation=man:bt-agent(1)
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
# DISABLED: Audio routing to prevent audio output to Pi
# The standard bt-agent accepts all service connections including audio profiles (A2DP, HSP, HFP)
# This would route audio output to the Raspberry Pi when devices connect.
# To re-enable audio routing, remove the --disable-audio flag (if supported) or use the custom agent instead.
# ExecStart=/usr/bin/bt-agent -c NoInputNoOutput
ExecStart=/bin/true
Restart=on-failure
RestartSec=5

# Don't restart if manually stopped
# on-failure means: restart on non-zero exit, killed by signal (except SIGHUP, SIGINT, SIGTERM, SIGPIPE)

[Install]
WantedBy=multi-user.target
