# presence-tracker.service
[Unit]
Description=IEEE Presence Tracker - Bluetooth Device Monitoring Service
After=bluetooth.target network-online.target
Wants=bluetooth.target

[Service]
Type=simple
User=ieee
WorkingDirectory=/home/ieee/Desktop/IEEE Presence Tracker
EnvironmentFile=/home/ieee/Desktop/IEEE Presence Tracker/.env
ExecStart=/home/ieee/.local/bin/uv run presence_tracker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
