# web-dashboard.service
[Unit]
Description=IEEE Presence Tracker Web Dashboard
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=ieee
Environment="PATH=/home/ieee/.local/bin:/usr/local/bin:/usr/bin:/bin"
WorkingDirectory=/home/ieee/Desktop/IEEE Presence Tracker
ExecStart=/home/ieee/.local/bin/uv run python web_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
